# meta-accelsius/meta-rcu/recipes-phosphor/gpio/phosphor-gpio-monitor/obmc-led-flashingON.service
[Unit]
Description=ID Button Pressed ON Service

[Service]
Restart=no
ExecStart=/usr/bin/env led-flashingON.sh
