# meta-accelsius/meta-rcu/recipes-phosphor/gpio/phosphor-gpio-monitor/obmc-led-flashingOFF.service
[Unit]
Description=ID Button Pressed OFF Service

[Service]
Restart=no
ExecStart=/usr/bin/env led-flashingOFF.sh
